# supervisor/cmit_faceswap.ini
[program:cmit_faceswap]
command=${项目根目录}/venv/bin/gunicorn -b 0.0.0.0:5050 -w 2 cmit_faceswap:app --timeout=180; supervisor启动命令
directory=${项目根目录}                                               ; 项目的文件夹路径
startsecs=0                                                                             ; 启动时间
stopwaitsecs=0                                                                          ; 终止等待时间
autostart=true                                                                         ; 是否自动启动
autorestart=true                                                                       ; 是否自动重启
stdout_logfile=${项目根目录}/logs/gunicorn.log                           ; log 日志
stderr_logfile=${项目根目录}/logs/gunicorn.err